# ansible/roles/sportspuff-app/templates/cache-warmer.timer.j2
[Unit]
Description=Timer for Sportspuff v6 Cache Warmer
Requires=cache-warmer.service

[Timer]
# Run every 20 seconds to keep scores cache warm (scores TTL is 30s)
# OnUnitActiveSec runs the timer 20 seconds after the service completes
OnUnitActiveSec=20s
# Also run on boot (with 30 second delay)
OnBootSec=30s
# If missed, run as soon as possible
Persistent=true

[Install]
WantedBy=timers.target

