# ansible/roles/sportspuff-app/templates/cache-warmer.timer.j2
[Unit]
Description=Timer for Sportspuff v6 Cache Warmer
Requires=cache-warmer.service

[Timer]
# During game hours (typically 6 PM - 2 AM), run every minute to keep scores fresh
# This covers most game times across timezones
OnCalendar=*-*-* 18:00..23:59:00
OnCalendar=*-*-* 00:00..02:00:00
# During off-hours, run every 5 minutes (schedules change less frequently)
OnCalendar=*-*-* 02:00..17:59:00/5
# Also run on boot (with 30 second delay)
OnBootSec=30s
# If missed, run as soon as possible
Persistent=true

[Install]
WantedBy=timers.target

