# ansible/roles/sportspuff-app/templates/nba-scores-updater.service.j2
[Unit]
Description=NBA Scores Updater Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_directory }}
Environment=PATH={{ venv_path }}/bin
ExecStart={{ venv_path }}/bin/python nba_scores_updater.py
StandardOutput=journal
StandardError=journal
# Don't restart on failure since this is a oneshot service
Restart=no
# Timeout after 5 minutes
TimeoutSec=300

[Install]
WantedBy=multi-user.target
