# ansible/roles/sportspuff-app/templates/nba-scores-updater.timer.j2
[Unit]
Description=NBA Scores Updater Timer
Requires=nba-scores-updater.service

[Timer]
# Run every 5 minutes during NBA season (roughly Oct-June)
OnCalendar=*:0/5
# Only run during NBA season months (Oct-June)
Persistent=true
# Start immediately if missed
AccuracySec=1min

[Install]
WantedBy=timers.target
